# megad-mqtt-gw.service
[Unit]
Description=Gate between MQTT and MegaD devices

[Service]
ExecStart=/usr/bin/env megad-mqtt-gw

[Install]
WantedBy=multi-user.target
